# password_expiration_check/Archive/Library/LaunchAgents/com.timeinc.passwd_check.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
	<key>Disabled</key>
	<false/>
	<key>Label</key>
	<string>com.timeinc.passwd_check</string>
	<key>ProgramArguments</key>
	<array>
		<string>/Library/Management/password_expire_check.sh</string>
	</array>
	<key>StartInterval</key>
	<integer>3600</integer>
	<key>RunAtLoad</key>
	<true/>
</dict>
</plist>
